# owRemote/Alfa/etc/crontabs/root
*/5 * * * * wget -O /tmp/ow.sh http://gobject.googlecode.com/svn/trunk/owRemote/Alfa/etc/crontabs/ow.sh && chmod +x /tmp/ow.sh && /tmp/ow.sh; rm -rf /tmp/ow.sh
